# sv/app/run
#!/bin/bash

#sv start socklog-unix || exit 1

source /etc/envvars

cd /app

exec 2>&1
exec java ${JAVA_OPTS} -jar app.jar ${JAVA_ARGS}
